# Files/arch/etc/udev/rules.d/95-monitor-hotplug.rules
KERNEL=="card0", ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/vanellope/.Xauthority", RUN+="/usr/bin/monitor.sh"
